# supervisor/apache2.conf
[program:apache2]
user=root
command=/usr/sbin/apache2ctl -DFOREGROUND
stdout_logfile=/var/log/supervisor/%(program_name).log
stderr_logfile=/var/log/supervisor/%(program_name)_error.log
